# layers/meta-elepi-userspace-odyssey/recipes-quests/hear-the-kernel-talking/files/systemd/honeypot.service
[Unit]
Description=Yummy yummy
After=systemd-journald.socket

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStartPre=modprobe kfrag
ExecStart=/usr/bin/honeypot.sh
ExecStopPost=rmmod kfrag

[Install]
WantedBy=multi-user.target