# layers/meta-elepi-userspace-odyssey/recipes-quests/hear-the-kernel-talking/files/systemd/honeypot.service
[Unit]
Description=Yummy yummy
After=run-postinsts.service

[Service]
ExecStart=/usr/bin/honeypot.sh

[Install]
WantedBy=multi-user.target